# deploy/roles/the_combine_app/templates/combine_backend.service.j2
[Unit]
Description=The Combine Backend service.

[Service]
User={{ combine_user }}
WorkingDirectory={{ combine_backend_dir }}/publish
{% for item in combine_backend_env_vars|dict2items %}
Environment="{{ item.key }}={{ item.value }}"
{% endfor %}
{% for item in combine_private_env_vars|dict2items %}
Environment="{{ item.key }}={{ item.value }}"
{% endfor %}
ExecStart=/usr/bin/dotnet {{ combine_backend_dir }}/publish/BackendFramework.dll
Restart=always
# Restart service after 10 seconds if the dotnet service crashes:
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier={{ combine_program_id }}

[Install]
WantedBy=multi-user.target
