# deploy/ansible/roles/support_tools/templates/display-cert-enddate.timer.j2
[Unit]
Description=Timer to trigger display-eth.service

[Timer]
OnUnitActiveSec={{ enddate_update_period }}
OnBootSec={{ enddate_update_period }}

[Install]
WantedBy=timers.target
