# deploy/ansible/roles/support_tools/templates/display-eth.timer.j2
[Unit]
Description=Display Ethernet Address on support tool display

[Timer]
OnUnitActiveSec={{ eth_update_period }}
OnBootSec={{ eth_update_period }}

[Install]
WantedBy=timers.target
