# deploy/ansible/roles/support_tools/templates/display-eth.service.j2
[Unit]
Description=Display ethernet address on external device
Wants=display-eth-addr.timer

[Service]
Type=oneshot
ExecStart={{ eth_update_program }}

[Install]
WantedBy=multi-user.target
